# root/etc/uci-defaults/luci-app-socat
#!/bin/sh

rm -rf /tmp/luci-*cache
exit 0
